# configdb3-deploy/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /lco/configdb3/configdb3-deploy/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
